# monerod.service
# systemd user service for monerod
# Install to ~/.config/systemd/user/ via: make install-services
# Set WALLET in ~/.config/p2pool/env
[Unit]
Description=Monero daemon (monerod)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/p2pool/env
ExecStart=monerod \
    --zmq-pub tcp://127.0.0.1:18083 \
    --out-peers 32 --in-peers 64 \
    --add-priority-node=p2pmd.xmrvsbeast.com:18080 \
    --add-priority-node=nodes.hashvault.pro:18080 \
    --enforce-dns-checkpointing \
    --enable-dns-blocklist \
    --prune-blockchain \
    --sync-pruned-blocks \
    --db-sync-mode=safe:sync \
    --data-dir %h/.bitmonero \
    --non-interactive
Restart=on-failure
RestartSec=30
TimeoutStopSec=120

[Install]
WantedBy=default.target
